# infra/hosts/lat1/systemd/finite-saas-runner.timer
[Unit]
Description=Run Finite agent creation runner

[Timer]
OnBootSec=30s
OnUnitInactiveSec=5s
AccuracySec=1s
Unit=finite-saas-runner.service

[Install]
WantedBy=timers.target
